# templates/mcewakeup.rules.j2
SUBSYSTEM=="usb", ATTRS{idVendor}=="1934", ATTRS{idProduct}=="5168" RUN+="/bin/sh -c 'echo enabled > /sys$env{DEVPATH}/../power/wakeup'"